# INSTALL/conf/dev_worker.ini
[program:dev-worker]
directory=/home/dev/reformWeb
environment=PATH="/home/dev/venv/bin:%(ENV_PATH)s"
command=/home/dev/venv/bin/python /home/dev/venv/bin/rq worker dev_queue --logging_level=INFO
autostart=true
autorestart=true
stderr_logfile=/var/log/reform/dev-worker.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stdout_logfile=/var/log/reform/dev-worker.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
user=dev
chown=dev:dev
