# INSTALL/conf/dev_reform.ini
[program:dev-reform]
directory=/home/dev/reformWeb
command=/home/dev/venv/bin/gunicorn app:app -b localhost:9000
autostart=true
autorestart=true
stderr_logfile=/var/log/reform/dev-reform.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stdout_logfile=/var/log/reform/dev-reform.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
user=dev
environment=REFORM_ENV="development"
