# INSTALL/conf/supervisor_reform.ini
[program:reform]
directory=/home/reform/reformWeb
command=/home/reform/venv/bin/gunicorn app:app -b localhost:8000
autostart=true
autorestart=true
stderr_logfile=/var/log/reform/reform.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stdout_logfile=/var/log/reform/reform.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
user=reform
